# src/buildroot/overlay/etc/sv/httpd/run
#!/bin/sh

exec httpd -f -h /var/www -u nobody
